# configs/z-manip-runtime-observer.service
[Unit]
Description=Z-Manipulation subscribe-only ROS2 runtime observer
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
Environment=ROS_DOMAIN_ID=20
ExecStart=%h/Z-Robotics-Lab/Z-Manipulation-Stack/scripts/runtime/go2w_runtime_observer.sh run
ExecStop=%h/Z-Robotics-Lab/Z-Manipulation-Stack/scripts/runtime/go2w_runtime_observer.sh stop
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
Nice=5
CPUWeight=50
IOWeight=50

[Install]
WantedBy=default.target
