# configs/z-manip-planning-workbench.service
[Unit]
Description=Z-Manip loopback planning-only workbench
After=default.target docker.service

[Service]
Type=simple
WorkingDirectory=%h/Z-Robotics-Lab/Z-Mobile-manip
ExecStart=/usr/bin/python3 %h/Z-Robotics-Lab/Z-Mobile-manip/scripts/runtime/go2w_planning_control.py --run-root %h/Z-Robotics-Lab/artifacts/go2w_real/planning_sessions --session-script %h/Z-Robotics-Lab/Z-Mobile-manip/scripts/runtime/go2w_planning_session.sh --home-script %h/Z-Robotics-Lab/Z-Mobile-manip/scripts/runtime/piper_home_remote.sh --grasp-script %h/Z-Robotics-Lab/Z-Mobile-manip/scripts/runtime/piper_full_grasp_remote.py --approach-script %h/Z-Robotics-Lab/Z-Mobile-manip/scripts/runtime/go2w_depth_servo.sh --approach-status %h/Z-Robotics-Lab/artifacts/go2w_real/latest/depth-servo.json --wrist-search-script %h/Z-Robotics-Lab/Z-Mobile-manip/scripts/runtime/piper_wrist_search_remote.sh --index %h/Z-Robotics-Lab/Z-Mobile-manip/web/debug_dashboard/index.html --runtime-state %h/Z-Robotics-Lab/artifacts/go2w_real/latest/runtime-observer.json --camera-image %h/Z-Robotics-Lab/artifacts/go2w_real/latest/camera-latest.jpg --depth-image %h/Z-Robotics-Lab/artifacts/go2w_real/latest/depth-latest.jpg --cloud-bin %h/Z-Robotics-Lab/artifacts/go2w_real/latest/cloud-latest.bin --port 8766
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
