# configs/z-manip-local-grounding.service
[Unit]
Description=Z-Mobile-Manip loopback-only YOLOE GPU grounding service
After=default.target docker.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f z-manip-local-grounding
ExecStart=/usr/bin/docker run --rm --name z-manip-local-grounding --gpus all --network host --read-only --tmpfs /tmp:rw,noexec,nosuid,size=128m -e MPLCONFIGDIR=/tmp/matplotlib z-mobile-manip-yoloe:latest
ExecStop=-/usr/bin/docker stop -t 5 z-manip-local-grounding
Restart=on-failure
RestartSec=2
TimeoutStartSec=120

[Install]
WantedBy=default.target
